# docker/supervisord.conf
; Supervises the processes that make up Wolf. Today that's PulseAudio + Wolf;
; drop in another [program:...] block to add more (e.g. Wolf Den).
;
; Env (PULSE_SERVER, WOLF_EMBED_PULSE, ...) is exported by startup.sh and
; inherited here. Program logs are forwarded to the container's stdout/stderr.

[supervisord]
nodaemon=true
user=root
; nodaemon already mirrors supervisord's own log to the console, so send the
; file copy to /dev/null to avoid duplicate lines in `docker logs`.
logfile=/dev/null
loglevel=info
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

; PulseAudio, started before Wolf and restarted if it ever dies at runtime.
; autostart is gated on WOLF_EMBED_PULSE so an external PULSE_SERVER is left alone.
; -n: skip the distro default.pa; auth-anonymous=1: let app containers connect
; without a cookie; exit-idle-time=-1: never shut down on its own.
[program:pulseaudio]
command=pulseaudio -n --daemonize=false --fail=false --exit-idle-time=-1 --load="module-native-protocol-unix auth-anonymous=1 socket=%(ENV_PULSE_SERVER)s" --load=module-always-sink --log-level=warn --log-target=stderr
autostart=%(ENV_WOLF_EMBED_PULSE)s
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; Wolf. When we manage PulseAudio, wait for its socket before starting so Wolf
; connects on the first try (the dependency, expressed without a fixed timeout).
[program:wolf]
command=bash -c 'if [ "$WOLF_EMBED_PULSE" = true ]; then echo "[wolf] waiting for PulseAudio socket $PULSE_SERVER"; until [ -S "$PULSE_SERVER" ]; do sleep 0.2; done; fi; exec /wolf/wolf'
autostart=true
autorestart=true
startretries=3
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
